# .config/wallpaper/schedule-update-wallpaper.service
[Unit]
Description=A job to update the wallpaper

[Service]
Type=oneshot
ExecStart=/home/%u/.config/wallpaper/update.sh